# internal_docs/01-strategic-planning/TEMPLATE-xnai-zram.service
[Unit]
Description=XNAi zRAM Setup - Phase 5A Memory Optimization
Documentation=https://wiki.archlinux.org/title/Zram
Before=swap.target
DefaultDependencies=no
After=systemd-modules-load.service
ConditionPathExists=/home/arcana-novai/Documents/xnai-foundation/scripts/xnai-zram-init.sh

[Service]
Type=oneshot
ExecStart=/bin/bash /home/arcana-novai/Documents/xnai-foundation/scripts/xnai-zram-init.sh
RemainAfterExit=yes

# Cleanup on stop
ExecStop=/usr/sbin/swapoff /dev/zram0
ExecStop=/usr/sbin/zramctl --reset /dev/zram0

StandardOutput=journal
StandardError=journal
SyslogIdentifier=xnai-zram

[Install]
WantedBy=swap.target
