# internal_docs/01-strategic-planning/TEMPLATE-xnai-zram.service
[Unit]
Description=XNAi zRAM Setup - Phase 5A Memory Optimization
Documentation=https://wiki.archlinux.org/title/Zram
Before=swap.target
DefaultDependencies=no
After=systemd-modules-load.service

[Service]
Type=oneshot
# Ensure wrapper exists (uses auto-sizing and nproc detection)
AssertPathExists=/home/arcana-novai/Documents/xnai-foundation/scripts/xnai-zram-init.sh
TimeoutStartSec=30
ExecStart=/bin/bash /home/arcana-novai/Documents/xnai-foundation/scripts/xnai-zram-init.sh
RemainAfterExit=yes
ExecStartPost=/bin/sh -c 'sleep 1; test -e /dev/zram0 || (journalctl -u xnai-zram.service -n 50; exit 1)'

# Cleanup on stop
ExecStop=/bin/swapoff /dev/zram0
ExecStop=/usr/bin/zramctl --reset /dev/zram0

StandardOutput=journal
StandardError=journal
SyslogIdentifier=xnai-zram

[Install]
WantedBy=swap.target
